# build/config-follower/runit/app/run
#!/bin/sh
node /app/bin/dist/app.js --config /app/conf/config-follower.json
