# build/config-follower/runit/app/run
#!/bin/sh
chpst -u :9999:9999 node /app/bin/dist/app.js --config /app/conf/config-follower.json
